# sw/src/overlay/etc/crontabs/root
*/1   *   *   *  *    /usr/bin/lua /usr/sbin/heatingd.lua