# src/etc/crontabs/root
*/1   *   *   *  *    /usr/bin/lua /root/heating.lua